# scripts/pipulse.service
[Unit]
Description=PiPulse — Raspberry Pi OLED Status
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pipulse
Restart=on-failure
# Requires membership in the `spi` and `gpio` groups.
# If running as root, these lines can be omitted.
SupplementaryGroups=spi gpio

[Install]
WantedBy=multi-user.target
